# linux/systemd/user/check-updates.service
[Unit]
Description=Check for dnf updates and notify

[Service]
Type=oneshot
ExecStart=%h/.dotfiles/scripts/check-updates.sh
