# linux/systemd/user/check-updates.timer
[Unit]
Description=Check for dnf updates daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
