# linux/etc/systemd/system/sshd.service.d/tailscale.conf
[Unit]
After=tailscaled.service
Wants=tailscaled.service
